# systemd-units/dropbox-folder-monitor.path
[Unit]
Description="Monitor dropbox folder changes"

[Path]
PathModified=/home/raw/MyDropbox
Unit=dropbox.service

[Install]
WantedBy=default.target
