# systemd-units/dropbox.service
[Unit]
Description=Dropbox Daemon sync

[Service]
ExecStart=/home/raw/.i3/sync_dropbox.sh
